# crontab

0 0 * * * root poetry run python /app/set_default_time.py
